# lcd_temp/airmsgd/CONTROL/postrm
#!/bin/sh
update-rc.d $D airmsgd remove
rm /usr/bin/airmsgd
